# psplash/initrd/scripts/init-bottom/end.psplash.service
[Unit]
Description=Stop Splash Screen
After=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill psplash

[Install]
WantedBy=getty.target
